# backend/papycast.service
[Unit]
Description=Papycast systemd service
After=network.target

[Service]
Type=simple
User=pi
WorkingDirectory=/opt/git/iia/papycast/backend
ExecStart=/opt/git/iia/papycast/backend/start.sh
Restart=always

[Install]
WantedBy=multi-user.target
